# infra/installer/linux/tpt-api.service
[Unit]
Description=TPT Online Video — API Service
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=tpt
Group=tpt
WorkingDirectory=/opt/tpt
ExecStart=/opt/tpt/bin/tpt-api
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
Environment=TPT_CONFIG=/opt/tpt/config/config.yaml

[Install]
WantedBy=multi-user.target