# infra/installer/linux/tpt-live.service
[Unit]
Description=TPT Online Video — Live Helper Service
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=tpt
Group=tpt
WorkingDirectory=/opt/tpt
ExecStart=/opt/tpt/bin/tpt-live
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
Environment=TPT_CONFIG=/opt/tpt/config/config.yaml

[Install]
WantedBy=multi-user.target