# infra/installer/linux/tpt-live.service
[Unit]
Description=TPT Online Video — Live Streaming Helper
After=network.target tpt-postgresql.service tpt-redis.service tpt-mediamtx.service
Requires=tpt-postgresql.service tpt-redis.service

[Service]
Type=simple
User=tpt
Group=tpt
WorkingDirectory=/opt/tpt
ExecStart=/opt/tpt/bin/tpt-live --config /opt/tpt/config/config.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
